# celery-worker.service
[Unit]
Description=Celery Worker for Inscription Service
After=network.target redis.service

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/root/registro_microservivio
Environment=PATH=/root/registro_microservivio/.venv/bin
EnvironmentFile=/root/registro_microservivio/.env
ExecStart=/root/registro_microservivio/.venv/bin/celery -A app.core.celery_app worker --loglevel=info --concurrency=4 --pidfile=/var/run/celery/worker.pid
ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
PIDFile=/var/run/celery/worker.pid

[Install]
WantedBy=multi-user.target